# deploy/rtm-workspace-sync.service
[Unit]
Description=Synchronize RTM developer workspace to Google Drive
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=rtmdeploy
Group=rtmdeploy
EnvironmentFile=/opt/rtm-app/.env
ExecStart=/usr/local/bin/rtm-sync-developer-workspace
